# systemd/openclaw-guard-bridge.timer
[Unit]
Description=Run OpenClaw Guard Bridge Tick every 2s

[Timer]
OnBootSec=10s
OnUnitActiveSec=2s
Unit=openclaw-guard-bridge.service

[Install]
WantedBy=timers.target
