# systemd/openclaw-cdp-watchdog.timer
[Unit]
Description=Run OpenClaw CDP watchdog periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=20s
Persistent=true

[Install]
WantedBy=timers.target
