# systemd/openclaw-cdp-watchdog.service
[Unit]
Description=OpenClaw CDP watchdog (restart webtop browser if CDP is down)
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/op-and-chloe/scripts/cdp-watchdog.sh
