# systemd/openclaw-guard-bridge.service
[Unit]
Description=OpenClaw Guard Bridge Tick
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/root/op-and-chloe/scripts/guard-bridge-tick.sh
